# deployment/pisad.service
[Unit]
Description=PISAD Beacon Detection Service
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=pisad
Group=pisad
WorkingDirectory=/home/pisad/projects/pisad
Environment="PATH=/home/pisad/.local/bin:/usr/local/bin:/usr/bin:/bin"
Environment="UV_SYSTEM_PYTHON=1"
Environment="PYTHONPATH=/home/pisad/projects/pisad/src"
ExecStart=/home/pisad/.local/bin/uv run uvicorn src.backend.core.app:app --host 0.0.0.0 --port 8080
Restart=on-failure
RestartSec=5
StartLimitBurst=3
StartLimitInterval=60
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
